# btt/crontab
0 * * * * ~/workspace/dotfiles/touchbar.sh
