# btt/crontab
0 * * * * ~/workspace/dotfiles/btt/touchbar.sh

